# packaging/root/opt/obsrvbl-ona/system/supervisord/ona-supervisord.conf
[supervisord]
logfile=/tmp/ona-supervisord.log
pidfile=/tmp/ona-supervisord.pid
nodaemon=false
user=obsrvbl_ona
logfile_maxbytes=1MB
logfile_backups=0

[program:ona-service]
command=/opt/obsrvbl-ona/system/supervisord/ona-service.sh
user=obsrvbl_ona
directory=/opt/obsrvbl-ona
autostart=true
autorestart=true
startretries=100
stopasgroup=true
redirect_stderr=true
stdout_logfile=/opt/obsrvbl-ona/logs/ona_service/ona-service.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=0
