# packaging/root/opt/obsrvbl-ona/system/upstart/obsrvbl-ona.conf
start on (filesystem and net-device-up IFACE!=lo) or (started JOB=network)

chdir /opt/obsrvbl-ona/
respawn

exec /usr/bin/sudo -u obsrvbl_ona /opt/obsrvbl-ona/system/supervisord/ona-supervisord.sh
